# docker/reporting/etc/supervisor/conf.d/daily_workhour_calculation.conf
[program:daily_workhour_calculation]
autostart = false
autorestart = false
command = php /srv/bin/console rabbitmq:consumer daily_workhour_calculation
startretries = 3