# docker/reporting/etc/supervisor/conf.d/reporting_employee.conf
[program:reporting_employee]
autostart = false
autorestart = false
command = php /srv/bin/console rabbitmq:consumer reporting_employee
startretries = 3