# docker/saga/etc/supervisor/conf.d/hrbase_employee.conf
[program:hrbase_employee]
autostart = false
autorestart = false
command = php /srv/bin/console rabbitmq:consumer hrbase_employee
startretries = 3