# docker/saga/etc/supervisor/conf.d/workhour_registration.conf
[program:workhour_registration]
autostart = false
autorestart = false
command = php /srv/bin/console rabbitmq:consumer workhour_registration
startretries = 3